# release/linux/fpm.before-remove.sh
systemctl stop detahardd.service
systemctl disable detahardd.service
