# release/linux/fpm.after-install.sh
systemctl enable detahardd.service
systemctl start detahardd.service
